# statusled.service
[Unit]
Description=statusled service

[Service]
Type=simple
ExecStart=python3 /usr/bin/statusled/serviceled.py

[Install]
WantedBy=multi-user.target